# manager.install/share/rustdesk.service.d/10-rdmanager-enable-socks.conf
[Service]
# pre configure proxy socks
ExecStartPre=/usr/sbin/lliurex-rdesk-manager socks-enable no-restart
